# usr/share/archboot/latest_installer/etc/profile.d/latest-installer.sh
cd /
echo "Waiting 10 seconds for getting an internet connection through dhcpcd..."
sleep 10
update-installer.sh -latest-install
