# usr/share/archboot/base/etc/systemd/system/detect-vconsole.path
# SPDX-License-Identifier: LGPL-2.1-or-later
# by Tobias Powalowski <tpowa@archlinux.org>

[Path]
PathModified=/sys/class/graphics/fbcon

[Install]
WantedBy=sysinit.target
