# etc/rc.local
#!/bin/sh
#
# /etc/rc.local: Local startup script.
#
# End of file
