# usr/share/archboot/net/etc/profile.d/elinks.sh
#!/bin/sh
export WWW_HOME="bit.ly/archboot"
export ELINKS_CONFDIR="/etc/elinks"
