# usr/share/archboot/net/etc/profile.d/elinks.sh
#!/bin/sh
export WWW_HOME="wiki.archlinux.org"
export ELINKS_CONFDIR="/etc/elinks"
