# usr/share/archboot/net/etc/profile.d/elinks.sh
#!/bin/sh
export WWW_HOME="https://archboot.com"
export ELINKS_CONFDIR="/etc/elinks"
