# usr/share/archboot/net/etc/profile.d/weechat.sh
#!/bin/sh
export WEECHAT_HOME="/etc/weechat"
