# usr/share/archboot/base/etc/profile.d/archboot-bash-aliases.sh
# SPDX-License-Identifier: GPL-2.0-only
# created by Tobias Powalowski <tpowa@archlinux.org>
# color output
alias ls='ls --color=auto'
alias grep='grep --color=auto'
# set neovim default editor
alias vi='nvim'
alias vim='nvim'
alias edit='nvim'
