# usr/lib/systemd/system/archboot-server-update.service
[Unit]
Description=Check, build and upload archboot images at night

[Service]
ExecStart=/usr/bin/archboot-server-update.sh

