# base/etc/rc.local
#!/bin/sh
#
# /etc/rc.local: Local multi-user startup script.
#
# End of file
