# usr/share/archboot/server/systemd/archboot.service
[Unit]
Description=Check, build and upload archboot images at night

[Service]
ExecStart=/home/tobias/Arch/iso/chroots/trigger-update.sh

